# ops/bissel-backup.timer
[Unit]
Description=Run remnant DB backup weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
