# ops/bissel-backup.timer
[Unit]
Description=Run remnant DB backup weekly

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target